# update-instance.conf
start on runlevel [2345]
stop on runlevel [!2345]

task

script
  su - ubuntu -c 'bash /home/ubuntu/dl-machine/scripts/install-deeplearning-libraries.sh reset http'
end script
